# conf/etc/s6/services/00-crond/run
#!/bin/bash

service-handler "CROD" "cron -f"